# deploy/crontab
00 9 * * 1-5 /home/ec2-user/canopy/canopy/src/core/RunAlgo.sh >/dev/null 2>&1